# config/files/common/usr/etc/profile.d/liquid-prompt-interactive.sh
[[ $- = *i* ]] && source /etc/liquidprompt/liquidprompt
